# archiso/airootfs/etc/systemd/system/synapse-kmod-build.service
[Unit]
Description=Build synapse_kmod for current kernel
# Only run if the module is not already built (any compression format)
ConditionPathExists=!/lib/modules/%v/extra/synapse_kmod.ko
ConditionPathExists=!/lib/modules/%v/extra/synapse_kmod.ko.zst
ConditionPathExists=!/lib/modules/%v/extra/synapse_kmod.ko.xz
ConditionPathExists=!/lib/modules/%v/extra/synapse_kmod.ko.gz
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/bin/synapse-kmod-build
RemainAfterExit=yes
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
