# archiso/airootfs/etc/systemd/system/synapd.service
[Unit]
Description=SynapseOS AI Daemon
After=network.target

[Service]
Type=simple
ExecStartPre=/bin/mkdir -p /run/synapd /var/lib/synapd
ExecStartPre=/bin/chmod 0755 /run/synapd
ExecStart=/usr/bin/synapd --foreground \
    --model /var/lib/synapd/models/synapse.gguf
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
