# archiso/airootfs/etc/systemd/system/synui.service
[Unit]
Description=SynapseOS Wayland Compositor
After=synapd.service seatd.service
Wants=synapd.service synui-foot.service
Requires=seatd.service

[Service]
Type=simple
User=root
TTYPath=/dev/tty1
StandardInput=tty
StandardOutput=tty
StandardError=journal
Environment=XDG_RUNTIME_DIR=/run/user/0
Environment=XDG_SESSION_TYPE=wayland
Environment=LIBSEAT_BACKEND=seatd
ExecStartPre=/bin/mkdir -p /run/user/0
ExecStart=/usr/bin/synui
Restart=on-failure
RestartSec=3
# If synui fails 3 times, stop trying and let getty take over
StartLimitBurst=3
StartLimitIntervalSec=30
ExecStopPost=-/usr/bin/systemctl start getty@tty1.service

[Install]
WantedBy=multi-user.target
