# archiso/airootfs/etc/systemd/system/var-lib-synapd-models.mount
[Unit]
Description=SynapseOS model storage (tmpfs)
Before=synapd.service

[Mount]
What=tmpfs
Where=/var/lib/synapd/models
Type=tmpfs
Options=mode=0755,size=6G

[Install]
WantedBy=multi-user.target
