# synguard/systemd/synguard.service
# /etc/systemd/system/synguard.service
# SynapseOS Security Monitor

[Unit]
Description=SynapseOS AI Security Monitor (synguard)
Documentation=https://synapseos.dev/docs/synguard
After=synapd.service local-fs.target
Requires=synapd.service
PartOf=synapd.service

[Service]
Type=notify
NotifyAccess=main
User=root
Group=synapse

ExecStartPre=/bin/mkdir -p /var/log/synguard /var/lib/synguard
ExecStart=/usr/bin/synguard \
    --foreground \
    --mode audit \
    --rules /etc/synguard/rules.d/

ExecReload=/bin/kill -HUP $MAINPID

Restart=always
RestartSec=3s
TimeoutStartSec=30s

# Journal output with priority tags
StandardOutput=journal
StandardError=journal
SyslogIdentifier=synguard

# Security hardening
ProtectProc=invisible
PrivateTmp=true
NoNewPrivileges=false

[Install]
WantedBy=multi-user.target
