# synnet/systemd/synnet.service
[Unit]
Description=SynapseOS AI Network Policy Daemon
After=synapd.service
Requires=synapd.service

[Service]
Type=simple
ExecStart=/usr/bin/synnet --foreground
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
