# synapd/systemd/synapd.service
# /etc/systemd/system/synapd.service
# SynapseOS AI Daemon — systemd service unit
#
# Core OS service. Starts before display manager and all
# synapse-aware services. Kernel-native inference always on.

[Unit]
Description=SynapseOS AI Daemon (synapd)
Documentation=https://synapseos.dev/docs/synapd
After=local-fs.target network.target
Before=display-manager.service synguard.service synnet.service
Wants=local-fs.target
DefaultDependencies=no

[Service]
Type=notify
NotifyAccess=main

# Root required for synapse_kmod sysfs interface
# Future: capability drop after init via setuid to 'synapse' user
User=root
Group=synapse

ExecStartPre=/bin/mkdir -p /run/synapd /var/lib/synapd /var/log/synapd
ExecStartPre=/bin/chmod 0750 /run/synapd
ExecStart=/usr/bin/synapd \
    --foreground \
    --model /var/lib/synapd/models/synapse-7b-q4_k_m.gguf \
    --socket /run/synapd/synapd.sock \
    --threads 4 \
    --context 4096

ExecReload=/bin/kill -HUP $MAINPID

# Always restart — synapd is a core OS service
Restart=always
RestartSec=2s
TimeoutStartSec=120s   # model loading can take a moment
TimeoutStopSec=30s

# Resource limits
LimitNOFILE=65536
LimitMEMLOCK=infinity   # needed for model mmap
LimitNICE=-20

# Security hardening
ProtectProc=invisible
PrivateTmp=true
NoNewPrivileges=false   # needs capabilities for sysfs

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=synapd

[Install]
WantedBy=multi-user.target
